# packages/peek-core-user/peek-core-user-5.1.6.tar.gz/peek_core_user-5.1.6/docker-ldap/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
loglevel=info

[program:samba]
command=/usr/sbin/samba --foreground --no-process-group
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/samba.err.log
stdout_logfile=/var/log/supervisor/samba.out.log
redirect_stderr=true