# packages/synerty-peek/synerty-peek-5.1.5.tar.gz/synerty_peek-5.1.5/ci/service/peek_office.service
[Unit]
Description=Peek Office Service
After=network.target

[Service]
LimitNOFILE=500000
Type=simple
User=peek
ExecStart=SED_PEEK_SOFTWARE_HOME/peek-platform/venv/bin/run_peek_office_service
Environment="SQLITE_TMPDIR=SED_PEEK_TMP_HOME"
WorkingDirectory=SED_PEEK_TMP_HOME
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target