# packages/synerty-peek/synerty-peek-5.1.5.tar.gz/synerty_peek-5.1.5/ci/service/peek_logic.service
[Unit]
Description=Peek Logic Service
After=network.target

[Service]
Type=simple
User=peek
ExecStart=SED_PEEK_SOFTWARE_HOME/peek-platform/venv/bin/run_peek_logic_service
Environment="SQLITE_TMPDIR=SED_PEEK_TMP_HOME"
Environment="LD_LIBRARY_PATH=SED_PEEK_PG_LD_LIBRARY_PATH"
WorkingDirectory=SED_PEEK_TMP_HOME
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target