# packages/getlino/getlino-25.10.0.tar.gz/getlino-25.10.0/getlino/templates/daphne_supervisor.conf
# generated by getlino

#[group:{{prjname}}]
#programs={{prjname}}_asgi,{{prjname}}_runworker,{{prjname}}_initiate_linod


[program:{{prjname}}_linod]
command=/bin/bash -c "source {{envdir}}/bin/activate && python manage.py linod --force"
#stopsignal=TERM
#killasgroup=true
priority=997
startsecs=5

# same for every program:
#stdout_logfile = {{project_dir}}/log/lino.log
redirect_stderr = true
environment = DJANGO_SETTINGS_MODULE={{django_settings_module}}
directory = {{project_dir}}
username = www-data
umask = 002

[fcgi-program:{{prjname}}_asgi]

# TCP socket used by Nginx backend upstream
# socket=unix://{{project_dir}}/nginx.sock
# socket_owner=root:www-data
# socket_mode=0660
# Multi-process daphne instances does not work with upstream unix sockets.
# See: https://github.com/django/daphne/issues/2871
socket=tcp://127.0.0.1:{{upstream_port}}

# Each process needs to have a separate socket file, so we use process_num
# Make sure to update "mysite.asgi" to match your project name
command=/bin/bash -c "source {{envdir}}/bin/activate && daphne -u /run/daphne/{{prjname}}_daphne%(process_num)d.sock --fd 0 --access-log - --proxy-headers asgi:application"

# Number of processes to startup, roughly the number of CPUs you have
numprocs={{num_procs}}

# Give each process a unique name so they can be told apart
process_name={{prjname}}_asgi%(process_num)d

# Automatically start and recover processes
autostart=true
autorestart=true
priority=998

# same for every program:
#stdout_logfile = {{project_dir}}/log/lino.log
redirect_stderr = true
environment = DJANGO_SETTINGS_MODULE={{django_settings_module}}
directory = {{project_dir}}
username = www-data
umask = 002
