# packages/synerty-peek/synerty-peek-5.1.6.tar.gz/synerty_peek-5.1.6/ci/service/peek_agent.service
[Unit]
Description=Peek Agent Service
After=network.target

[Service]
Type=simple
User=peek
ExecStart=SED_PEEK_SOFTWARE_HOME/peek-platform/venv/bin/run_peek_agent_service
Environment="SQLITE_TMPDIR=SED_PEEK_TMP_HOME"
Environment="LD_LIBRARY_PATH=SED_PEEK_ORACLE_HOME"
Environment="ORACLE_HOME=SED_PEEK_ORACLE_HOME"
Environment="SQLITE_TMPDIR=SED_PEEK_TMP_HOME"
WorkingDirectory=SED_PEEK_TMP_HOME
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target